# puppet/zulip_ops/files/supervisor/conf.d/grafana.conf
[program:grafana]
command=/opt/grafana/bin/grafana-server -config /etc/grafana/grafana.ini
directory=/opt/grafana
priority=15
autostart=true
autorestart=true
user=grafana
